# deploy/paopao-launch-radar.service
[Unit]
Description=PaoPao Launch Radar Web API
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=__PROJECT_DIR__
EnvironmentFile=-__PROJECT_DIR__/.env.oi
Environment=PYTHONUNBUFFERED=1
ExecStart=__PYTHON_BIN__ -u main.py web --web-host 127.0.0.1 --web-port 18090 --web-mode real
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
